# contrib/init/straksd.conf
description "Straks Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env STRAKSD_BIN="/usr/bin/straksd"
env STRAKSD_USER="straks"
env STRAKSD_GROUP="straks"
env STRAKSD_PIDDIR="/var/run/straksd"
# upstart can't handle variables constructed with other variables
env STRAKSD_PIDFILE="/var/run/straksd/straksd.pid"
env STRAKSD_CONFIGFILE="/etc/straks/straks.conf"
env STRAKSD_DATADIR="/var/lib/straksd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # straksd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$STRAKSD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run straksd."
        echo "The setting must appear in $STRAKSD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Straks Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$STRAKSD_PIDDIR"
    chmod 0755 "$STRAKSD_PIDDIR"
    chown $STRAKSD_USER:$STRAKSD_GROUP "$STRAKSD_PIDDIR"
    chown $STRAKSD_USER:$STRAKSD_GROUP "$STRAKSD_CONFIGFILE"
    chmod 0660 "$STRAKSD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$STRAKSD_PIDFILE" \
    --chuid $STRAKSD_USER:$STRAKSD_GROUP \
    --exec "$STRAKSD_BIN" \
    -- \
    -pid="$STRAKSD_PIDFILE" \
    -conf="$STRAKSD_CONFIGFILE" \
    -datadir="$STRAKSD_DATADIR" \
    -disablewallet \
    -daemon

